# roles/certificates/templates/cert-renew.timer
[Unit]
Description=Certificate renewal timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
