# install/playbooks/roles/privoxy/templates/update-rules.cron
# privoxy rules auto update
# Update rules every day at 3am
0 3 * * * /usr/local/sbin/update-rules.sh > /var/log/privoxy/rules-update.log 2>&1
